# com.mattia.kavita-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd"\>
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mattia.kavita-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>while true; do /opt/homebrew/bin/fswatch -1 "/Users/mattiacolombo/Calibre Library/" && /Users/mattiacolombo/.kavita-sync.sh; done</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/kavita-sync.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/kavita-sync.log</string>
</dict>
</plist>
